# config/security/systemd/weekly-blog-agentic.service.d/10-ssh-identity.conf.example
[Service]
Environment="GIT_SSH_COMMAND=ssh -F /home/user/.config/casonk.github.io/ssh/config"
